# init/etc/inittab
::restart:/sbin/reboot

::respawn:/etc/initshell.sh
