# init/etc/inittab
#::sysinit:/etc/init.d/rcS

::restart:/sbin/reboot
